# templates/tomcat.init.erb
#!/bin/bash
#
# Init file for Tomcat server
#
# chkconfig: 2345 55 25
# description: Tomcat server
#

# Source function library.
. /etc/init.d/functions

TOMCAT_USER="<%= tomcat_user %>"
CATALINA_HOME="<%= basedir %>"
PID_FILE="<%= basedir %>/temp/tomcat.pid"

PROG="tomcat"

help() {
    echo "$0 INVALID SYNTAX ::  Use:"
	echo "$0 {stop | start | restart | status}"
}

start() {
	echo -n "Starting ${PROG}: "
	su --login --command "${CATALINA_HOME}/bin/startup.sh" ${TOMCAT_USER}
	ps -u ${TOMCAT_USER} --no-headers -o pid,cmd | grep "${CATALINA_HOME}" | awk '{ print $1 }' > ${PID_FILE};
	echo "done.";
}

stop() {
    echo "Stopping ${PROG}: "

    su --login --command "${CATALINA_HOME}/bin/shutdown.sh" ${TOMCAT_USER}

    # Use the pid file if it exists.
    if [[ -f ${PID_FILE} ]]; then
	    for  i in 1 2 3; do
            if [[ -d /proc/`cat ${PID_FILE}` ]]; then
                kill `cat ${PID_FILE}`
                sleep 2
            fi
        done
        rm -f ${PID_FILE}
    else
        # The tomcat procs for tomcat cle user
        PIDS=`ps -u ${TOMCAT_USER} --no-headers -o pid,cmd | grep "${CATALINA_HOME}" | awk '{ print $1 }'`
        for PID in ${PIDS}; do
            for  i in 1 2 3; do
                if [[ -d /proc/${PID} ]]; then
                    kill ${PID}
                    sleep 2
                fi
            done
            
            if [[ -d /proc/${PID} ]]; then
                kill -9 ${PID}
            fi
        done
    fi
    echo 'done'
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
esac

exit 0
