# templates/tomcat.init.erb
#!/bin/bash
#
# Init file for Tomcat server
#
# chkconfig: 2345 55 25
# description: Tomcat server
#

# Source function library.
. /etc/init.d/functions

TOMCAT_USER="<%= tomcat_user %>"
CATALINA_HOME="<%= basedir %>"

PROG="tomcat"

help() {
    echo "$0 INVALID SYNTAX ::  Use:"
	echo "$0 {stop | start | restart | status}"
}

start() {
	echo -n "Starting ${PROG}: "
	su --login --command "${CATALINA_HOME}/bin/startup.sh" ${TOMCAT_USER}
	ps -u ${TOMCAT_USER} --no-headers -o pid,cmd | grep "${CATALINA_HOME}" | awk '{ print $1 }' > ${INIT_PID};
	echo "done.";
}

function killcle() {
    MAX_TRIES=3
    TRIES=1
    PID=${1}

    echo "killing ${PROG} with pid=${PID}"

    while [[ $TRIES < $MAX_TRIES ]]; do
        if [[ -d /proc/${1} ]]; then
            su --login --command "kill ${PID}" ${TOMCAT_USER}
            sleep 5
        fi
        TRIES=$[TRIES+1]
    done

    if [[ -d /proc/${PID} ]]; then
        kill -9 ${PID}
    fi
}

stop() {
    echo "Stopping ${PROG}: "

    su --login --command "{$CATALINA_HOME}/bin/shutdown.sh" ${TOMCAT_USER}

    # Use the pid file if it exists.
    if [[ -f ${INIT_PID} ]]; then
	    for PID in `cat ${INIT_PID}`; do
            killcle ${PID}
	    done
        rm -f ${INIT_PID}
    else
        # The CLE procs for the cle user
        PIDS=`ps -u ${TOMCAT_USER} --no-headers -o pid,cmd | grep "${CATALINA_HOME}" | awk '{ print $1 }'`
        for PID in ${PIDS}; do
            if [[ $? -eq 0 ]]; then
                killcle ${PID}
            fi
        done
    fi
    echo 'done'
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
esac

exit 0
