# modules/oae/templates/solr.erb
#!/bin/bash
#
# chkconfig: 2345 01 99
# description: Solr cluster node.

# Global, editable VARS
SOLR_JAVA_OPTS=" -Dsolr.solr.home=<%= solr_app %> -Djetty.port=8983 "
SAKAIOAE_USER="<%= scope.lookupvar('oae::params::user') %>"
SOLR_HOME=<%= solr_app %>
JAVA_BIN="/usr/bin/java"

INIT_PID="/var/run/oae.solr.pid"

PROG="solr"

help() {
echo "$0 INVALID SYNTAX ::  Use:"
echo "$0 {stop | start | restart | status}"
}

start() {
echo -n "Starting $PROG: "
su --login --command "(cd ${SOLR_HOME} && ${JAVA_BIN} ${SOLR_JAVA_OPTS} -jar ${SOLR_HOME}/start.jar > /dev/null 2>&1 &)" ${SAKAIOAE_USER}
echo $$ > ${INIT_PID};
echo "done.";
}

stop() {
echo "Stopping $PROG: "
PROCESS=`ps auxwww | grep "^${SAKAIOAE_USER}" | perl -pe 's/^\S*\s*([0-9]+).*\n/\1 /'`
if [ "${PROCESS}" != "" ] ; then
  su -c "kill ${PROCESS}" ${SAKAIOAE_USER}
  sleep 5;
  PROCESS=`ps auxwww | grep "^${SAKAIOAE_USER}" | perl -pe 's/^\S*\s*([0-9]+).*\n/\1 /'`
  if [ "${PROCESS}" != "" ] ; then
    su -c "kill -9 ${PROCESS}" ${SAKAIOAE_USER}
    sleep 5;
    PROCESS=`ps auxwww | grep "^${SAKAIOAE_USER}" | perl -pe 's/^\S*\s*([0-9]+).*\n/\1 /'`
    if [ "${PROCESS}" != "" ] ; then
      echo "failed.";
      exit 1;
    fi
  fi
fi
rm -f ${INIT_PID};
echo "done.";
}

restart() {
stop
start
}

status() {
if [ -f ${INIT_PID} ]; then
  echo "It seems that the $SOLR app is still running.  If you have issued a stop command already, please wait." ;
  exit 0;
else
  echo "SakaiOAE is not running" ;
  exit 2;
fi
}

case $1 in
    "start")
        start;
        ;;
    "stop")
        stop;
        ;;
    "status")
        status;
        ;;
    "help")
        help;
        ;;
    *)
        help
        ;;
esac

exit 0;
